# rpm/micro-ecommerce/src/e-commerce-cart.service
[Unit]  
Description=e-commerce-cart

[Service]  
ExecStart=/opt/e-commerce/e-commerce-cart -Dspring.config.location=/etc/e-commerce/cart/application.yml
ExecStop=/bin/kill -TERM $MAINPID  
ExecReload=/bin/kill -HUP $MAINPID  
Type=simple
Restart=on-failure
KillMode=control-group

[Install]  
WantedBy=multi-user.target  
